# systemd/cafe-coffee-generate.service
[Unit]
Description=Cafe Today Coffee — daily insight self-generation (Claude API)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/root/92cafe/cafe-today-coffee
EnvironmentFile=/root/92cafe/cafe-today-coffee/.env
ExecStart=/bin/bash /root/92cafe/cafe-today-coffee/scripts/generate.sh
StandardOutput=journal
StandardError=journal
TimeoutStartSec=20min
